# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
directory=/fedcloud-dashboard
command=gunicorn --bind 0.0.0.0:8000 dashboard.main:app

[program:cron]
command=cron -f
